# init.d/genoclaim.service
[Unit]
Description=genoclaim

[Service]
ExecStart=/usr/local/bin/genoclaim -k -d /var/lib/genoclaim --logdir /var/log/genoclaim -p

[Install]
WantedBy=multi-user.target
